# go-wild-dns.service
[Unit]
Description=Go Wild DNS

[Service]
ExecStart=/root/go/bin/go-wild-dns
Restart=on-abort

[Install]
WantedBy=multi-user.target

